# infra/service
#!/bin/sh
# nuc/apps/burpee/service
# Run as root on NUC — installs daemon(8) rc.d script and env template into burpee jail
# Run once after nuc/apps/burpee/jail
set -e

SCRIPT_DIR="$(cd "$(dirname "$0")" && pwd)"
. "${SCRIPT_DIR}/config"

JROOT="/jails/${JAIL}"

log() { echo "==> [${APP}] $*"; }

log "Creating runtime directories"
# /var/run is tmpfs-like on FreeBSD — create via rc.d or here; we use a
# persistent path under /var/run inside the jail and create it on every start
# via the rc.d script instead, so the pidfile survives reboots correctly.
jexec "$JAIL" mkdir -p "/var/run/${JAIL}" "/var/log/${JAIL}" "/var/db/${JAIL}"
jexec "$JAIL" chown -R "${JAIL}:${JAIL}" \
  "/var/run/${JAIL}" "/var/log/${JAIL}" "/var/db/${JAIL}"

log "Writing rc.d service"
cat >"${JROOT}/usr/local/etc/rc.d/${APP}" <<EOF
#!/bin/sh
#
# PROVIDE: ${APP}
# REQUIRE: NETWORKING
# KEYWORD: shutdown

. /etc/rc.subr

name="${APP}"
rcvar="${APP}_enable"

pidfile="/var/run/${JAIL}/${APP}.pid"
logfile="/var/log/${JAIL}/${APP}.log"

load_rc_config \$name
: \${${APP}_enable:="NO"}
: \${${APP}_user:="${JAIL}"}

start_cmd="${APP}_start"
stop_cmd="${APP}_stop"
status_cmd="${APP}_status"

${APP}_start() {
    mkdir -p /var/run/${JAIL}
    set -a && . /etc/${JAIL}/env && set +a
    /usr/sbin/daemon \
        -S \
        -r \
        -R 5 \
        -u \${${APP}_user} \
        -P \${pidfile} \
        -o \${logfile} \
        /opt/${JAIL}/bin/${APP} start
}

${APP}_stop() {
    [ -f "\${pidfile}" ] \
        && kill "\$(cat \${pidfile})" && rm -f "\${pidfile}" \
        || echo "\${name} not running"
}

${APP}_status() {
    [ -f "\${pidfile}" ] && kill -0 "\$(cat \${pidfile})" 2>/dev/null \
        && echo "\${name} running as PID \$(cat \${pidfile})" \
        || { echo "\${name} not running"; return 1; }
}

run_rc_command "\$1"
EOF

chmod +x "${JROOT}/usr/local/etc/rc.d/${APP}"

grep -q "${APP}_enable" "${JROOT}/etc/rc.conf" ||
  echo "${APP}_enable=\"YES\"" >>"${JROOT}/etc/rc.conf"

log "Writing env template"
[ -f "${JROOT}/etc/${JAIL}/env" ] && {
  log "Env file already exists — skipping. Edit manually if needed."
} || cat >"${JROOT}/etc/${JAIL}/env" <<EOF
PHX_SERVER=true
PHX_HOST=${PHX_HOST}
PORT=${PORT}
DATABASE_PATH=/var/db/${JAIL}/${APP}.db
POOL_SIZE=5
SECRET_KEY_BASE=REPLACE_ME
RELEASE_COOKIE=REPLACE_ME
EOF

chmod 600 "${JROOT}/etc/${JAIL}/env"

log ""
log "==> Service installed."
log "Fill in ${JROOT}/etc/${JAIL}/env before first start:"
log "  - SECRET_KEY_BASE  (run: mix phx.gen.secret)"
log "  - RELEASE_COOKIE   (run: mix phx.gen.secret 32)"
log "Then: sh infra/deploy aktersnurra@nuc"
