# deploy/freebsd/video-server
#!/usr/bin/env bash
# nuc/apps/burpee/video-server
# Sets up a local Caddy instance inside the burpee jail to serve protected
# videos. Only accepts requests carrying the correct X-Video-Secret header.
# Run as root on NUC after nuc/apps/burpee/jail.
#
# Requires: VIDEO_SECRET (must match VIDEO_SECRET used in vps/02-caddy)
set -euo pipefail

SCRIPT_DIR="$(cd "$(dirname "$0")" && pwd)"
. "${SCRIPT_DIR}/config"

VIDEO_SECRET="${VIDEO_SECRET:?set VIDEO_SECRET — must match the value used in vps/02-caddy}"
JROOT="/jails/${JAIL}"
VIDEO_DIR="/var/lib/${APP}/videos"

log() { echo "==> [video-server] $*"; }

log "Installing Caddy in jail"
jexec "$JAIL" pkg install -y caddy

log "Writing Caddyfile"
mkdir -p "${JROOT}/etc/caddy"
cat >"${JROOT}/etc/caddy/videos.Caddyfile" <<EOF
{
    # No ACME — this server is only reachable over Tailscale from the VPS
    auto_https off
}

:${VIDEO_PORT} {
    # Reject any request that doesn't carry the correct shared secret.
    # This prevents the NUC from serving videos to anything other than
    # the VPS Caddy instance.
    @unauthorized not header X-Video-Secret "${VIDEO_SECRET}"
    respond @unauthorized 403

    root * ${VIDEO_DIR}
    file_server
}
EOF

SVC="${JAIL}_videos"

log "Writing rc.d service"
cat >"${JROOT}/usr/local/etc/rc.d/${SVC}" <<EOF
#!/bin/sh
# PROVIDE: ${SVC}
# REQUIRE: NETWORKING
# KEYWORD: shutdown

. /etc/rc.subr

name="${SVC}"
rcvar="${SVC}_enable"
pidfile="/var/run/${JAIL}/${SVC}.pid"
logfile="/var/log/${JAIL}/${SVC}.log"

load_rc_config \$name
: \${${SVC}_enable:="NO"}
: \${${SVC}_user:="${JAIL}"}

start_cmd="${SVC}_start"
stop_cmd="${SVC}_stop"
status_cmd="${SVC}_status"

${SVC}_start() {
    /usr/sbin/daemon \
        -S \
        -r \
        -R 5 \
        -u \${${SVC}_user} \
        -P \${pidfile} \
        -o \${logfile} \
        /usr/local/bin/caddy run --config /etc/caddy/videos.Caddyfile
}

${SVC}_stop() {
    [ -f "\${pidfile}" ] \
        && kill "\$(cat \${pidfile})" && rm -f "\${pidfile}" \
        || echo "\${name} not running"
}

${SVC}_status() {
    [ -f "\${pidfile}" ] && kill -0 "\$(cat \${pidfile})" 2>/dev/null \
        && echo "\${name} running as PID \$(cat \${pidfile})" \
        || { echo "\${name} not running"; return 1; }
}

run_rc_command "\$1"
EOF

chmod +x "${JROOT}/usr/local/etc/rc.d/${SVC}"

grep -q "${SVC}_enable" "${JROOT}/etc/rc.conf" ||
  echo "${SVC}_enable=\"YES\"" >>"${JROOT}/etc/rc.conf"

log "Ensuring video directory exists with correct ownership"
jexec "$JAIL" mkdir -p "$VIDEO_DIR"
jexec "$JAIL" chown -R "${JAIL}:${JAIL}" "$VIDEO_DIR"

log "Starting service"
jexec "$JAIL" service "${SVC}" start

log ""
log "==> Video file server running on jail port ${VIDEO_PORT}."
log "VPS Caddy will proxy /protected-videos/ to ${JAIL_IP}:${VIDEO_PORT} with X-Video-Secret."
log ""
log "Also add this pf rule on the NUC (nuc/01-base) once VPS_TS_IP is known:"
log "  pass in on \$ts_if proto tcp from <VPS_TS_IP> to ${JAIL_IP} port ${VIDEO_PORT} keep state"
